# ckan/files/ckan-supervisor.conf
{% from "ckan/map.jinja" import ckan with context -%}
[program:ckan]
command={{ ckan.venv_path }}/bin/gunicorn --paste {{ [ckan.confdir, ckan.conffile]|join('/') }} --bind :5000 --workers=2
user={{ ckan.ckan_user }}
redirect_stderr=true
stdout_logfile={{ ckan.ckan_home }}/ckan-gunicorn.log
autostart=false
autorestart=true
startsecs=10
