# ckan/files/ckan-supervisor.conf
{% from "ckan/map.jinja" import ckan with context -%}
[program:ckan]
command={{ ckan.venv_path }}/bin/paster serve {{ [ckan.confdir, ckan.conffile]|join('/') }}
user={{ ckan.ckan_user }}
numprocs=1
stdout_logfile={{ ckan.ckan_home }}/paster.log
stderr_logfile={{ ckan.ckan_home }}/paster.log
autostart=false
autorestart=true
startsecs=10
