# ckan/files/datapusher-supervisor.conf
{% from "ckan/map.jinja" import ckan with context -%}
{% set main = [ckan.src_dir, 'datapusher', 'datapusher',
               'main.py']|join('/') -%}
{% set settings = [ckan.src_dir, 'datapusher', 'deployment',
                   'datapusher_settings.py']|join('/') -%}
[program:datapusher]
command={{ ckan.venv_path }}/bin/python {{ main }} {{ settings }}
#command={{ ckan.venv_path }}/bin/gunicorn wsgi:app -b 127.0.0.1:8800 -c {{ settings }}
user={{ ckan.ckan_user }}
numprocs=1
stdout_logfile={{ ckan.ckan_home }}/datapusher.log
stderr_logfile={{ ckan.ckan_home }}/datapusher.log
autostart=true
autorestart=true
