# ckan/files/redis-supervisor.conf
{% from "ckan/map.jinja" import ckan with context %}
[program:redis]
command=redis-server --dir {{ ckan.ckan_home }} --dbfilename redis.db
user={{ ckan.ckan_user }}
numprocs=1
stdout_logfile={{ ckan.ckan_home }}/redis.log
stderr_logfile={{ ckan.ckan_home }}/redis.log
autostart=true
autorestart=true
startsecs=10
priority=900
