# ckan/extensions/files/archive-server.conf
{% from "ckan/map.jinja" import ckan with context -%}
[program:archive-server]
command=/usr/bin/python -m SimpleHTTPServer 9876
user={{ ckan.ckan_user }}
directory={{ ckan.extensions.archiver.options.get('ckanext-archiver.archive_dir',
                                                  '/tmp/archive') }}
stdout_logfile={{ [ckan.ckan_home, 'archiver-server.log']|join('/') }}
stderr_logfile={{ [ckan.ckan_home, 'archiver-server.log']|join('/') }}
autostart=false
autorestart=true
