# ckan/extensions/files/supervisor-harvest.conf
{% from "ckan/map.jinja" import ckan with context %}
[program:ckan_gather_consumer]
command={{ ckan.venv_path }}/bin/paster --plugin=ckanext-harvest harvester gather_consumer --config={{ [ckan.confdir, ckan.conffile]|join('/') }}
user={{ ckan.ckan_user }}
directory={{ [ckan.src_dir, 'ckan']|join('/') }}
numprocs=1
stdout_logfile={{ [ckan.ckan_home, 'gather_consumer.log']|join('/') }}
stderr_logfile={{ [ckan.ckan_home, 'gather_consumer.log']|join('/') }}
autostart = true
autorestart = true
startsecs=10

[program:ckan_fetch_consumer]
command={{ ckan.venv_path }}/bin/paster --plugin=ckanext-harvest harvester fetch_consumer --config={{ [ckan.confdir, ckan.conffile]|join('/') }}
user={{ ckan.ckan_user }}
directory={{ [ckan.src_dir, 'ckan']|join('/') }}
numprocs=1
stdout_logfile={{ [ckan.ckan_home, 'fetch_consumer.log']|join('/') }}
stderr_logfile={{ [ckan.ckan_home, 'fetch_consumer.log']|join('/') }}
autostart = true
autorestart = true
startsecs=10
