# ckan/extensions/files/supervisor-archiver.conf
{% from "ckan/map.jinja" import ckan with context %}
[program:archiver_priority]
command={{ ckan.venv_path }}/bin/paster --plugin=ckanext-archiver celeryd2 run priority --config={{ [ckan.confdir, ckan.conffile]|join('/') }}
user={{ ckan.ckan_user }}
numprocs=1
stdout_logfile={{ [ckan.ckan_home, 'archiver_priority.log']|join('/') }}
stderr_logfile={{ [ckan.ckan_home, 'archiver_priority.log']|join('/') }}
autostart = true
autorestart = true
startsecs=10

[program:archiver_bulk]
command={{ ckan.venv_path }}/bin/paster --plugin=ckanext-archiver celeryd2 run bulk --config={{ [ckan.confdir, ckan.conffile]|join('/') }}
user={{ ckan.ckan_user }}
numprocs=1
stdout_logfile={{ [ckan.ckan_home, 'archiver_bulk.log']|join('/') }}
stderr_logfile={{ [ckan.ckan_home, 'archiver_bulk.log']|join('/') }}
autostart = true
autorestart = true
startsecs=10
