# ckan/extensions/files/supervisor-harvest.conf
{% from "ckan/map.jinja" import ckan with context %}
[program:ckan_gather_consumer]
command={{ ckan.venv_path }}/bin/paster --plugin=ckanext-harvest harvester gather_consumer --config={{ [ckan.confdir, ckan.conffile]|join('/') }}
user={{ ckan.ckan_user }}
directory={{ [ckan.src_dir, 'ckan']|join('/') }}
numprocs=1
stdout_logfile={{ [ckan.ckan_home, 'gather_consumer.log']|join('/') }}
stderr_logfile={{ [ckan.ckan_home, 'gather_consumer.log']|join('/') }}
autostart = true
autorestart = true
startsecs=10
environment=CKAN_SITE_URL="{{ ckan.site_url }}",CKAN_SQLALCHEMY_URL="postgresql://{{ ckan.db_user }}:{{ ckan.db_password }}@{{ ckan.db_host }}/{{ ckan.db_name }}",CKAN_DATASTORE_WRITE_URL="postgresql://{{ ckan.db_user }}:{{ ckan.db_password }}@{{ ckan.db_host }}/{{ ckan.datastore_db_name }}",CKAN_DATASTORE_READ_URL="postgresql://{{ ckan.datastore_db_user }}:{{ ckan.datastore_db_password }}@{{ ckan.db_host }}/{{ ckan.datastore_db_name }}",CKAN_SITE_URL="{{ ckan.site_url }}",CKAN_STORAGE_PATH="{{ ckan.ckan_home }}/storage",CKAN_SOLR_URL="http://localhost:8983/solr"

[program:ckan_fetch_consumer]
command={{ ckan.venv_path }}/bin/paster --plugin=ckanext-harvest harvester fetch_consumer --config={{ [ckan.confdir, ckan.conffile]|join('/') }}
user={{ ckan.ckan_user }}
directory={{ [ckan.src_dir, 'ckan']|join('/') }}
numprocs=1
stdout_logfile={{ [ckan.ckan_home, 'fetch_consumer.log']|join('/') }}
stderr_logfile={{ [ckan.ckan_home, 'fetch_consumer.log']|join('/') }}
autostart = true
autorestart = true
startsecs=10
environment=CKAN_SITE_URL="{{ ckan.site_url }}",CKAN_SQLALCHEMY_URL="postgresql://{{ ckan.db_user }}:{{ ckan.db_password }}@{{ ckan.db_host }}/{{ ckan.db_name }}",CKAN_DATASTORE_WRITE_URL="postgresql://{{ ckan.db_user }}:{{ ckan.db_password }}@{{ ckan.db_host }}/{{ ckan.datastore_db_name }}",CKAN_DATASTORE_READ_URL="postgresql://{{ ckan.datastore_db_user }}:{{ ckan.datastore_db_password }}@{{ ckan.db_host }}/{{ ckan.datastore_db_name }}",CKAN_SITE_URL="{{ ckan.site_url }}",CKAN_STORAGE_PATH="{{ ckan.ckan_home }}/storage",CKAN_SOLR_URL="http://localhost:8983/solr"
