# solr/files/solr-supervisor.conf
{% from "ckan/map.jinja" import ckan,solr with context %}
[program:solr]
command = java -jar start.jar
user={{ ckan.ckan_user }}
directory={{ solr.home }}
numprocs=1
stdout_logfile={{ ckan.ckan_home }}/solr.log
stderr_logfile={{ ckan.ckan_home }}/solr.log
autostart=true
autorestart=true
